# homedir/Documents/Scripts/crontab
0 11 * * * bash /home/carneca/Documents/Scripts/daily_backup.sh
0 16 * * 2 bash /home/carneca/Documents/Scripts/weekly_backup.sh
