# homedir/Documents/Scripts/crontab
0 11 * * * bash /home/carneca/Documents/Scripts/daily_backup.sh
